# deploy/systemd/idena-session-recorder.service
[Unit]
Description=PoHW Idena forward session recorder
After=idena.service
Wants=idena.service
RequiresMountsFor=/mnt/ssd

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/mnt/ssd/p2pool/pohw_idena_rpc
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /mnt/ssd/p2pool/pohw_idena_rpc/idena_session_recorder.py --db /mnt/ssd/pohw-p2pool/idena-session-recorder/session_recorder.sqlite3 --rpc-url http://127.0.0.1:9009 --api-key-file /mnt/ssd/idena/idena-data/api.key --poll-interval 10 --max-blocks-per-pass 500 run
Restart=always
RestartSec=30
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadOnlyPaths=/mnt/ssd/p2pool /mnt/ssd/idena/idena-data
ReadWritePaths=/mnt/ssd/pohw-p2pool/idena-session-recorder
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
