# deploy/systemd/idena-public-state-export.service
[Unit]
Description=Export synchronized Idena public state and push it to the Pi
After=network-online.target idena-bootstrap.service
Wants=network-online.target
RequiresMountsFor=/srv/idena /srv/idena-return-export /var/lib/idena-return /etc/idena-return

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/idena-return/export.env
ExecStart=/usr/local/sbin/idena-public-state-export-push
ExecStopPost=/usr/local/sbin/idena-public-state-export-push --recover-source
TimeoutStartSec=2h
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/srv/idena /etc/idena-return
ReadWritePaths=/srv/idena-return-export /var/lib/idena-return /run/lock
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER
MemoryDenyWriteExecute=true
