# deploy/systemd/idena-hetzner-modern.service
[Unit]
Description=Idena original-chain node (modern Hetzner runtime)
Documentation=https://www.idena.io/guide
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/srv/idena
ConditionPathIsMountPoint=/srv/idena
StartLimitIntervalSec=5min
StartLimitBurst=3

[Service]
Type=simple
User=idena-modern
Group=idena-modern
WorkingDirectory=/srv/idena
Environment=HOME=/srv/idena
Environment=XDG_CONFIG_HOME=/srv/idena/.config
ExecStart=/usr/local/libexec/idena-node-compat-v6 --config=/etc/idena-modern/config.json --datadir=/srv/idena --fast=false
Restart=on-failure
RestartSec=30
TimeoutStopSec=5min
LimitNOFILE=65536
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
ReadOnlyPaths=/etc/idena-modern/config.json /usr/local/libexec/idena-node-compat-v6
ReadWritePaths=/srv/idena
InaccessiblePaths=/srv/idena-original-relay

[Install]
WantedBy=multi-user.target
