# deploy/systemd/idena-public-state-import.path
[Unit]
Description=Watch for a completed Idena public-state return transfer

[Path]
PathExists=/var/lib/idena-return-inbox/current/READY
PathExists=/var/lib/idena-return-state/in-progress.json
Unit=idena-public-state-import.service

[Install]
WantedBy=multi-user.target
