# deploy/systemd/pohw-health-status.timer
[Unit]
Description=Periodic PoHW host health status

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=pohw-health-status.service

[Install]
WantedBy=timers.target
