# deploy/systemd/pohw-network-watchdog.service
[Unit]
Description=PoHW network self-recovery watchdog
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/
EnvironmentFile=-/etc/pohw/p2pool.env
ExecStart=/usr/bin/env POHW_NETWORK_WATCHDOG_STATE_DIR=/var/lib/pohw/network-watchdog /usr/local/libexec/pohw/pohw-network-watchdog.sh
TimeoutStartSec=45
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_NET_RAW CAP_SYS_BOOT
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_NET_RAW CAP_SYS_BOOT
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
UMask=0077
ReadWritePaths=/var/lib/pohw/network-watchdog
