# deploy/systemd/pohw-idena-priority-guard.service
[Unit]
Description=PoHW Idena validation priority guard
After=network-online.target idena.service
Wants=network-online.target idena.service
RequiresMountsFor=/mnt/ssd

[Service]
Type=oneshot
WorkingDirectory=/mnt/ssd/p2pool
EnvironmentFile=-/etc/pohw/p2pool.env
ExecStart=/usr/bin/python3 /mnt/ssd/p2pool/scripts/pohw-idena-priority-guard.py
TimeoutStartSec=30
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_DAC_READ_SEARCH
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
UMask=0077
ReadOnlyPaths=/mnt/ssd/p2pool /mnt/ssd/idena/idena-data
ReadWritePaths=/mnt/ssd/pohw-p2pool
