# deploy/systemd/pohw-mining-adapter-sdcard.service
[Unit]
Description=PoHW P2Pool local Stratum mining adapter (SD-card runtime)
Wants=network-online.target pohw-gossip-mesh.service
After=network-online.target pohw-gossip-mesh.service pohw-fork-chain-node.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/p2pool
EnvironmentFile=-/etc/pohw/p2pool.env
ExecStart=/opt/p2pool/scripts/pohw-run-mining-adapter.sh
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
UMask=0077
ReadOnlyPaths=/opt/p2pool /etc/pohw
ReadWritePaths=/var/lib/pohw-p2pool

[Install]
WantedBy=multi-user.target
