# deploy/systemd/pohw-idena-snapshot-sdcard.service
[Unit]
Description=PoHW Idena snapshot when the local node is fully synced (SD-card runtime)
Wants=idena.service idena-reward-indexer.service
After=idena.service idena-reward-indexer.service network-online.target
RequiresMountsFor=/opt/p2pool /var/lib/idena /var/lib/pohw-p2pool

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/p2pool
EnvironmentFile=-/etc/pohw/p2pool.env
Environment=IDENA_API_KEY_FILE=/var/lib/idena/api.key
Environment=IDENA_REWARD_LEDGER_DB=/var/lib/pohw-p2pool/rewards/reward_ledger.sqlite3
Environment=IDENA_REWARD_INDEXER_SCRIPT=/opt/p2pool/pohw_idena_rpc/idena_reward_indexer.py
Environment=POHW_SNAPSHOT_DIR=/var/lib/pohw-p2pool/snapshots
ExecStart=/opt/p2pool/scripts/pohw-snapshot-if-synced.sh
TimeoutStartSec=15min
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
UMask=0077
ReadOnlyPaths=/opt/p2pool /var/lib/idena
ReadWritePaths=/var/lib/pohw-p2pool
