# deploy/systemd/pohw-dashboard-api-cookie-watch@.path
[Unit]
Description=Watch Bitcoin RPC cookie for PoHW dashboard API reload (%I)
RequiresMountsFor=/%I

[Path]
PathChanged=/%I
PathModified=/%I
Unit=pohw-dashboard-api-cookie-watch.service

[Install]
WantedBy=multi-user.target
