# deploy/systemd/pohw-dashboard-api-cookie-watch.path
[Unit]
Description=Watch Bitcoin RPC cookie for PoHW dashboard API reload
RequiresMountsFor=/mnt/ssd

[Path]
PathChanged=/mnt/ssd/bitcoin/bitcoin-core-mainnet/.cookie
PathModified=/mnt/ssd/bitcoin/bitcoin-core-mainnet/.cookie
Unit=pohw-dashboard-api-cookie-watch.service

[Install]
WantedBy=multi-user.target
